# systemd/screeninvader2.service
[Unit]
Description=ScreenInvader 2.0 Player + Webserver
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=screeninvader
Group=screeninvader

WorkingDirectory=/opt/screeninvader2

# mpv auf TTY3 ausführen – TTY1 bleibt für Terminal nutzbar
TTYPath=/dev/tty3
StandardInput=tty
StandardOutput=tty
StandardError=tty

# kleine Verzögerung, damit beim Boot kurz ein Terminal sichtbar bleibt
ExecStartPre=/bin/sleep 3

# Python-Server im venv starten (mpv wird von server.py / player.py gesteuert)
ExecStart=/opt/screeninvader2/.venv/bin/python server.py

Restart=always
RestartSec=2

# Zugriff auf TTY-Konfiguration (für mpv auf separatem TTY)
AmbientCapabilities=CAP_SYS_TTY_CONFIG

[Install]
WantedBy=multi-user.target
